# CSCU_Proto3/FileSystemFiles_Common/nandflash/etc/cron.d/root
0 4 * * * timer_upload.sh &
50 3 * * * timer_delete.sh &
